# systemd/unicorn.service
[Unit]
Description=UnicornSystem Service
[Service]
ExecStart=/opt/unicorn/unicorn_system
